# service/gamma-energy.service
[Unit]
Description=Web server that serves voltage data generated by gamma pulses.
Wants=network-online.target
After=network-online.target

[Service]
Restart=on-failure
RestartSec=5
Type=simple
WorkingDirectory=/root/energy-histogram/server/
ExecStart=/usr/local/bin/gunicorn 'app:create_app()' --bind 0.0.0.0:5000

[Install]
WantedBy=multi-user.target
