# BeagleBone/Service/GammaEnergy.service
[Unit]
Description=Redis-based script for reading voltages generated by gamma pulses.
Wants=network-online.target
After=network-online.target

[Service]
Restart=on-failure
RestartSec=5
Type=simple
WorkingDirectory=/root/energy-histogram/BeagleBone
ExecStart=/bin/bash /root/energy-histogram/BeagleBone/run.sh

[Install]
WantedBy=multi-user.target